# crontab
*/5 * * * * /app/dynamic-dns/cloudflare.sh
